# 95-HDMI2.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/opt/bin/hotplug_HDMI2.sh"
